# deploy/disk-hygiene-watchdog.service
[Unit]
Description=Disk-Hygiene Watchdog — Auto-Prune + Alarm bei vollem Disk
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-%h/.config/shadowops-watchdog.env
ExecStart=%h/shadowops-bot/scripts/disk-hygiene-watchdog.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=disk-hygiene-watchdog
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
