# deploy/guildscout-watchdog.timer
[Unit]
Description=GuildScout Watchdog Timer — alle 5 Minuten (mit Jitter)
Documentation=https://github.com/Commandershadow9/shadowops-bot

[Timer]
# Erster Lauf 4 Minuten nach Boot (versetzt zu shadowops:2 und zerodox:3)
OnBootSec=4min
OnUnitActiveSec=5min
RandomizedDelaySec=45s
AccuracySec=15s
Persistent=false

[Install]
WantedBy=timers.target
