# deploy/security-freshness-watchdog.service
[Unit]
Description=Security-Team Freshness Watchdog — letzter erfolgreicher sec_jobs-Lauf < 26h (DB completed_at)
Documentation=https://github.com/Commandershadow9/shadowops-bot
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/home/cmdshadow/.config/shadowops-watchdog.env
Environment=WATCHDOG_SERVICE_NAME=security-freshness
Environment=WATCHDOG_MODE=pg-freshness
Environment=WATCHDOG_PG_CONTAINER=guildscout-postgres
Environment=WATCHDOG_PG_USER=security_analyst
Environment=WATCHDOG_PG_DB=security_analyst
Environment=WATCHDOG_MAX_AGE_HOURS=26
Environment="WATCHDOG_PG_QUERY=SELECT COALESCE(EXTRACT(EPOCH FROM (NOW()-MAX(completed_at)))/3600, 999999) FROM sec_jobs WHERE status IN ('ok','partial')"
ExecStart=/home/cmdshadow/shadowops-bot/scripts/service-watchdog.sh
SuccessExitStatus=0 1
StandardOutput=journal
StandardError=journal
