# deploy/shadowops-bot.service
[Unit]
Description=ShadowOps Security Discord Bot - Active Security Guardian
Documentation=https://github.com/commandershadow/shadowops-bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cmdshadow
Group=cmdshadow
WorkingDirectory=/home/cmdshadow/shadowops-bot

# Cleanup stale PID file and stop orphaned process (if any)
ExecStartPre=/bin/bash -c 'PID_FILE=/home/cmdshadow/shadowops-bot/.bot.pid; if [ -f "$PID_FILE" ]; then pid=$(cat "$PID_FILE"); if [ -n "$pid" ] && kill -0 "$pid" 2>/dev/null; then kill "$pid"; sleep 2; if kill -0 "$pid" 2>/dev/null; then kill -9 "$pid"; fi; fi; rm -f "$PID_FILE"; fi'

# Python mit Bot-Script (uv Virtual Environment)
ExecStart=/home/cmdshadow/shadowops-bot/.venv/bin/python3 /home/cmdshadow/shadowops-bot/src/bot.py

# Graceful Shutdown: SIGTERM → close() mit per-Component Timeouts
# Security Engine allein braucht bis zu 20s, Discord-Close bis zu 15s
KillMode=mixed
TimeoutStopSec=120
SuccessExitStatus=143 SIGTERM

# Restart-Policy
Restart=on-abnormal
RestartSec=30
StartLimitIntervalSec=300
StartLimitBurst=3

# Limits
LimitNOFILE=65536
MemoryMax=3G

# Logs
StandardOutput=journal
StandardError=journal
SyslogIdentifier=shadowops-bot

# Security
# WICHTIG: false nötig für sudo cscli/fail2ban-client Zugriff
NoNewPrivileges=false
PrivateTmp=true

# Umgebungsvariablen
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONUTF8=1"
# Security Analyst DB — entweder hier ODER in config.yaml unter security_analyst.database_dsn
# Environment="SECURITY_ANALYST_DB_URL=postgresql://security_analyst:PASSWORD@127.0.0.1:5433/security_analyst"
# WICHTIG: CLAUDECODE darf NICHT gesetzt sein (verhindert "nested session" Fehler bei Claude CLI)
UnsetEnvironment=CLAUDECODE
# PATH fuer Codex und Claude CLI
Environment="PATH=/home/cmdshadow/.nvm/versions/node/v25.2.1/bin:/home/cmdshadow/.local/bin:/usr/local/bin:/usr/bin:/bin"
# XDG_RUNTIME_DIR für systemctl --user Zugriff auf Agent-Services
Environment="XDG_RUNTIME_DIR=/run/user/1000"

[Install]
WantedBy=multi-user.target
