# deploy/zerodox-watchdog.service
[Unit]
Description=ZERODOX Health Watchdog — externer Health-Check + Discord-Alert
Documentation=https://github.com/Commandershadow9/shadowops-bot
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/home/cmdshadow/.config/shadowops-watchdog.env
Environment=WATCHDOG_SERVICE_NAME=zerodox
Environment=WATCHDOG_HEALTH_URL=https://zerodox.de/api/health
Environment=WATCHDOG_REQUIRE_BOT_READY=0
# Webhook-URL kommt aus dem EnvironmentFile (gleicher wie shadowops-watchdog).
# Falls dort als SHADOWOPS_WATCHDOG_WEBHOOK gesetzt: Backward-Compat-Fallback
# im Script greift automatisch.
ExecStart=/home/cmdshadow/shadowops-bot/scripts/service-watchdog.sh
SuccessExitStatus=0 1
StandardOutput=journal
StandardError=journal
